# scripts/org.jleechan.user-scope-backup.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.jleechan.user-scope-backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bash</string>
    <string>@HOME@/projects_other/user_scope/scripts/backup-home.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>REQUIRE_DROPBOX_ONLY</key>
    <string>1</string>
    <key>ALLOW_GIT_BACKUP_SYNC</key>
    <string>0</string>
    <key>ALLOW_GIT_BACKUP_PUSH</key>
    <string>0</string>
    <key>TOKEN_HEALTH_PRECHECK</key>
    <string>1</string>
    <key>TOKEN_HEALTH_PRECHECK_FAIL_ON_FAILURE</key>
    <string>1</string>
    <key>GIT_BACKUP_PUSH_REMOTE</key>
    <string>origin</string>
    <key>GIT_BACKUP_PUSH_BRANCH</key>
    <string>main</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>7200</integer>
  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/user-scope-backup.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/user-scope-backup.launchd.err.log</string>
</dict>
</plist>
