# automation/com.jleechanorg.pr-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job identification -->
    <key>Label</key>
    <string>com.jleechanorg.pr-automation</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jleechan/projects/worktree_worker2/automation/jleechanorg_pr_monitor.py</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/jleechan/projects/worktree_worker2</string>

    <!-- Schedule: Every 10 minutes -->
    <key>StartInterval</key>
    <integer>600</integer>

    <!-- Auto-start on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep the job alive if it exits -->
    <key>KeepAlive</key>
    <false/>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/jleechan/.pyenv/shims:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/jleechan</string>
        <key>USER</key>
        <string>jleechan</string>
        <!-- GitHub CLI authentication -->
        <key>GITHUB_TOKEN</key>
        <string>$GITHUB_TOKEN</string>
        <!-- Optional: Custom workspace location -->
        <key>PR_AUTOMATION_WORKSPACE</key>
        <string>/Users/jleechan/tmp/jleechanorg-pr-workspaces</string>
    </dict>

    <!-- Standard output and error paths -->
    <key>StandardOutPath</key>
    <string>/Users/jleechan/Library/Logs/worldarchitect-automation/jleechanorg-launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jleechan/Library/Logs/worldarchitect-automation/jleechanorg-launchd.err</string>

    <!-- Process limits -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Throttle settings (prevent rapid restart) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Only run when user is logged in -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>NumberOfProcesses</key>
        <integer>100</integer>
    </dict>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Exit timeout -->
    <key>ExitTimeOut</key>
    <integer>3600</integer>
</dict>
</plist>
